# ttrss-updater/ttrss-updater.service
[Unit]
Description=TT-RSS Updater
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker exec ttrss_app php /var/www/update.php --daemon
User=root
